# skills/integration-automation/wterm/launchd/com.minoan.wterm-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.minoan.wterm-server</string>
    <key>ProgramArguments</key>
    <array>
        <string>__NODE__</string>
        <string>__HOME__/daimones/wterm-server/node_modules/.bin/tsx</string>
        <string>__HOME__/daimones/wterm-server/src/server.ts</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__HOME__/daimones/wterm-server</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__HOME__/.claude/logs/wterm-server.out.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.claude/logs/wterm-server.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:__HOME__/.local/bin</string>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>
</dict>
</plist>
